# route.cron
16 06 * * *  sh /home/cjw1115/route/route_add.sh >> /home/cjw1115/route/log_add.txt
